# rootfs/etc/s6-overlay/s6-rc.d/30-create-user-and-group/up
/etc/s6-overlay/s6-rc.d/30-create-user-and-group/run
