# rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/bin/bash
set -e
ME=sshd

entrypoint_log() {
	echo "$ME: $*"
}

entrypoint_log "INFO: Starting sshd server..."
/usr/sbin/sshd -e -D
