# rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/execlineb -P
/docker-entrypoint.sh /usr/sbin/sshd -e -D
